# kiosk/pi-gen-stage/02-kiosk-service/files/famcal-ota.timer
[Unit]
Description=FamCal OTA check every 6 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
