# kiosk/pi-gen-stage/02-kiosk-service/files/famcal-ota.service
[Unit]
Description=FamCal OTA Update Check

[Service]
Type=oneshot
ExecStart=/home/kiosk/famcal/ota-update.sh
